# config/download.conf
[program:download]
command = /data/code/web/bin/gunicorn_start.sh                        ; Command to start app
user = hunch                                                          ; User to run as
stdout_logfile = /data/code/web/logs/gunicorn_supervisor.log   	      ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8                       ; Set UTF-8 as default encoding
